# systemd/meshtastic-bridge-watchdog.service
[Unit]
Description=Meshtastic Bridge Watchdog

[Service]
Type=oneshot
ExecStart=/opt/meshtastic-matrix-bridge/watchdog.sh
